# duo/systemd/duo-watch-displays.service
# Reference template — the home-manager module (modules/zenbook-duo)
# generates the real unit with the correct repo path. For non-Nix use:
# copy to ~/.config/systemd/user/, fix ExecStart, then
#   systemctl --user daemon-reload && systemctl --user enable --now duo-watch-displays
[Unit]
Description=zenduo: keep the panel layout matching the keyboard dock state
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.dotfiles/duo/bin/duo watch-displays
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
